# board/recalbox/fsoverlay/etc/udev/rules.d/99-openjvs.rules
SUBSYSTEM=="tty", KERNEL=="ttyUSB0", ACTION=="add", RUN+="/usr/bin/python -m py_openjvs"