# projects/recalbox-hardware/case/installers/gpi2/assets/92-gpicase2-audio.rules
DRIVER=="snd-usb-audio", RUN+="/usr/bin/gpicase-audio"
